# pkgs/kde/caelestia-src/systemd/qs-kwin-bridge.service
[Unit]
Description=Quickshell KDE Window Bridge
PartOf=graphical-session.target
After=plasma-kwin_wayland.service

[Service]
ExecStart=/usr/bin/env python3 %h/.local/bin/qs-kwin-bridge.py
Restart=always
RestartSec=3
StartLimitIntervalSec=0

[Install]
WantedBy=graphical-session.target
